# fleet-quota.timer
[Unit]
Description=Refresh cached Agent Fleet quota every fifteen minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=1min
Persistent=true

[Install]
WantedBy=timers.target
